# Linux (debian) server/etc/systemd/system/backup@client1.timer
[Unit]
Description=Daily backup for client1

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target